# ops/deploy/ikimon_v2_production_media_smoke.timer
[Unit]
Description=Run ikimon production media smoke daily

[Timer]
OnCalendar=*-*-* 04:17:00
RandomizedDelaySec=20m
Persistent=true
Unit=ikimon-v2-production-media-smoke.service

[Install]
WantedBy=timers.target
