# ops/deploy/ikimon_v2_media_worker.service
[Unit]
Description=ikimon platform_v2 media processing worker
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/ikimon.life/repo
ExecStart=/var/www/ikimon.life/repo/ops/cron/run_v2_node_script.sh processMediaProcessingJobs --limit=10 --stale-seconds=900
Nice=5
