# ops/deploy/ikimon_v2_media_worker.timer
[Unit]
Description=Run ikimon platform_v2 media processing worker

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true
Unit=ikimon-v2-media-worker.service

[Install]
WantedBy=timers.target
