# ops/deploy/ikimon_v2_blue.service
[Unit]
Description=ikimon platform_v2 blue API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/ikimon.life/runtime/blue
EnvironmentFile=/etc/ikimon/production-v2.env
Environment=NODE_ENV=production
Environment=PORT=3201
ExecStart=/usr/bin/env node dist/server.js
Restart=always
RestartSec=5
UMask=0002
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
