# ops/deploy/ikimon_v2_guide_environment.service
[Unit]
Description=ikimon platform_v2 guide environment mesh and regional hypothesis refresh
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/ikimon.life/repo
ExecStart=/var/www/ikimon.life/repo/ops/cron/run_v2_node_script.sh runGuideEnvironmentPostDeploy --date=today --trigger=timer --hypothesis-limit=500 --eval-limit=2000
Nice=10
