# ops/deploy/ikimon_v2_alert_delivery_worker.timer
[Unit]
Description=Run ikimon platform_v2 alert delivery worker

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=ikimon-v2-alert-delivery-worker.service

[Install]
WantedBy=timers.target
