# platform_v2/ops/systemd/ikimon-cold-cache-invalidate.timer
[Unit]
Description=Trigger ikimon-cold-cache-invalidate every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=ikimon-cold-cache-invalidate.service

[Install]
WantedBy=timers.target
